# ops/devnet/ippan-export-dataset.service
[Unit]
Description=IPPAN Devnet Dataset Export (D-GBDT telemetry)
After=network-online.target ippan-node.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
TimeoutStartSec=900
WorkingDirectory=/root/IPPAN
StandardOutput=journal
StandardError=journal
ExecStart=/usr/local/lib/ippan/export-dataset.sh


