# deploy/ippan-node.service
[Unit]
Description=IPPAN Blockchain Node
Documentation=https://github.com/ippan/ippan
After=network.target
Wants=network.target

[Service]
Type=simple
User=ippan
Group=ippan
WorkingDirectory=/opt/ippan
ExecStart=/opt/ippan/bin/ippan-node
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ippan-node

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/ippan
ReadWritePaths=/var/log/ippan
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=32768

# Environment
EnvironmentFile=/etc/ippan/ippan.env

# Kill settings
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target